# config/ansible/roles/vagrant-user/files/profile.d/opt-dirs-in-path.sh
export PATH="$PATH:/opt/consul/bin:/opt/nomad/bin"